# daemon/edo-daemon.service
[Unit]
Description=edo-daemon (CTFd infra sidecar)
After=network-online.target docker.service wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
# Read EDO_ALLOWED_UID + tuneables from a root-owned drop-in.
EnvironmentFile=/etc/edo/daemon.env
# edo_daemon.py imports the edo_core/ package as a sibling directory, so the
# whole daemon/ folder must be installed together — not just this one file.
# ExecStart uses the venv's own interpreter (NOT /usr/bin/python3) because
# modern Debian/Ubuntu block `pip install` into system Python (PEP 668) and
# a plain venv is invisible to a root-owned systemd unit unless pointed at
# directly. See README: Installation, section 1.
WorkingDirectory=/opt/edo/daemon
ExecStart=/opt/edo/venv/bin/python3 /opt/edo/daemon/edo_daemon.py
Restart=on-failure
RestartSec=3
# Give us just enough — we shell out to iptables/wg/docker so we need
# NET_ADMIN and the ability to reach /var/run/docker.sock. CAP_CHOWN is
# required even though we run as root: CapabilityBoundingSet restricts
# which capabilities a process can use AT ALL, regardless of UID, and
# os.chown() on the RPC socket is gated specifically by CAP_CHOWN — NOT by
# CAP_DAC_OVERRIDE, which only bypasses permission-bit (rwx) checks, not
# ownership changes. Without it, chown() fails with EPERM even for UID 0.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_OVERRIDE CAP_CHOWN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_CHOWN
NoNewPrivileges=yes
ProtectSystem=strict
# Dash prefix = "OK if missing at boot" — /etc/wireguard doesn't exist until
# wireguard-tools is installed, and /run/edo/ is created by the daemon itself.
ReadWritePaths=-/run/edo -/var/lib/edo -/etc/wireguard
# The daemon creates its own runtime directory with the right ownership.
# RuntimeDirectoryPreserve=restart is critical: without it, systemd deletes
# and recreates /run/edo as a NEW directory on every restart (deploys,
# crash-recovery, manual restarts). Docker's bind mount of /run/edo into
# the CTFd container is anchored to the directory INSTANCE that existed
# when the container started — a fresh instance orphans that mount, and
# the socket silently vanishes from inside the container until the CTFd
# container is also recreated. `restart` preserves the directory across
# Stop+Start cycles (still cleaned up on a genuine `systemctl stop`).
RuntimeDirectory=edo
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=restart
ProtectHome=yes
# ProtectHome=yes makes /home, /root, and /run/user entirely invisible to
# the daemon (docs: "mostly equivalent to setting [them] in
# InaccessiblePaths="). Consequence: every challenge's build_path MUST
# live outside those three trees, or the daemon will report "no
# Dockerfile" even though the file is genuinely there — e.g. /opt/... or
# /srv/..., not anything under /root or /home. (ReadOnlyPaths= and
# BindReadOnlyPaths= were both tried as exceptions for a path under
# /root and neither actually worked against this InaccessiblePaths-style
# masking on systemd 255 — verified via nsenter showing no submount ever
# got created — so don't reach for those; move the content instead.)
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
