# daemon/edo-daemon.service
[Unit]
Description=edo-daemon (CTFd infra sidecar)
After=network-online.target docker.service wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
# Read EDO_ALLOWED_UID + tuneables from a root-owned drop-in.
EnvironmentFile=/etc/edo/daemon.env
# edo_daemon.py imports the edo_core/ package as a sibling directory, so the
# whole daemon/ folder must be installed together — not just this one file.
# ExecStart uses the venv's own interpreter (NOT /usr/bin/python3) because
# modern Debian/Ubuntu block `pip install` into system Python (PEP 668) and
# a plain venv is invisible to a root-owned systemd unit unless pointed at
# directly. See README: Installation, section 1.
WorkingDirectory=/opt/edo/daemon
ExecStart=/opt/edo/venv/bin/python3 /opt/edo/daemon/edo_daemon.py
Restart=on-failure
RestartSec=3
# Give us just enough — we shell out to iptables/wg/docker so we need
# NET_ADMIN and the ability to reach /var/run/docker.sock.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=yes
ProtectSystem=strict
# Dash prefix = "OK if missing at boot" — /etc/wireguard doesn't exist until
# wireguard-tools is installed, and /run/edo/ is created by the daemon itself.
ReadWritePaths=-/run/edo -/var/lib/edo -/etc/wireguard
# The daemon creates its own runtime directory with the right ownership.
RuntimeDirectory=edo
RuntimeDirectoryMode=0755
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
